# app/crontab

# * * * * * echo "$(date)" >> /tmp/cron.log

* * * * * php artisan top:send-mail-command 1

